# private_dot_config/systemd/user/pr-digest.service
[Unit]
Description=Generate desktop notification with GitHub PR focus state

[Service]
Type=oneshot
ExecStart=%h/.config/my-scripts/bin/pr-digest.sh
